# conf/debrepo.crontab
* * * * * debrepo [ -x /home/debrepo/scripts/cleanup.sh ] && /home/debrepo/scripts/cleanup.sh

